# s6/update-daemon/run
#!/bin/sh

exec s6-setuidgid nobody /usr/bin/php /var/www/update_daemon2.php
