# src/s6/update-daemon/run
#!/bin/sh

exec s6-setuidgid nobody php82 /var/www/update_daemon2.php
