# src/s6/update-daemon/run
#!/bin/sh

exec s6-setuidgid www-data php /var/www/update_daemon2.php
