# s6/php-fpm/run
#!/bin/sh

exec /usr/bin/php-fpm --nodaemonize
