# src/s6/php-fpm/run
#!/bin/sh

PHP_SUFFIX=${PHP_SUFFIX:-83}

exec "/usr/sbin/php-fpm${PHP_SUFFIX}" --nodaemonize
