# src/s6/php-fpm/run
#!/bin/sh

exec /usr/sbin/php-fpm81 --nodaemonize
