# src/s6/php-fpm/run
#!/bin/sh

exec /usr/local/sbin/php-fpm --nodaemonize
