# src/s6/php-fpm/run
#!/bin/sh

exec /usr/sbin/php-fpm82 --nodaemonize
